# rootfs/etc/services.d/fail2ban/run
#!/usr/bin/execlineb -P

fail2ban-server -f
